# config/com.user.plexmonitor.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.user.plexmonitor</string>
    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/plex/plex_monitor.sh</string>
    </array>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>StandardOutPath</key>
    <string>/Users/plex/Library/Logs/plex_monitor.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/plex/Library/Logs/plex_monitor_error.log</string>
    <key>UserName</key>
    <string>plex</string>
</dict>
</plist>

